# hk-smzjg-components.service
[Unit]
Description=SMZJG 所有组件统一启动服务
After=network.target
Wants=network.target

[Service]
User=root
Type=simple

ExecStart=/opt/hk/smzjg/script/component/start-components.sh
ExecStop=/opt/hk/smzjg/script/component/stop-components.sh all

# 设置 Restart 策略用于崩溃时自动拉起
Restart=on-failure
RestartSec=5s

# 日志输出
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

